# docker/overlay/etc/periodic/hourly/preload
#!/bin/sh

set -euo pipefail


#======================================================================================================================
# Call preload script.
#======================================================================================================================

bf-write "Calling preload script." preload
ob preload
